# data/pkg/skel/etc/service/cgrates/run
#!/bin/sh
echo Firing up CGRateS engine
exec setuidgid cgrates /usr/bin/cgr-engine -config /etc/cgrates/cgrates.cfg

