# data/pkg/skel/etc/service/cgrates/run
#!/bin/sh
echo Running CGRateS service
exec /usr/bin/cgr-engine -config /etc/cgrates/cgrates.cfg

