# debian/cgrates.cgrates-backup.timer
[Unit]
Description=Daily CGRateS backup

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=900

[Install]
WantedBy=timers.target
